# ccminer.service
[Unit]
Description=ccminer Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/root/start.sh
User=root
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
